# .config/autostart/ibus-start.sh.desktop
[Desktop Entry]
Exec=/home/shima/ibus-start.sh
Icon=
Name=ibus-start.sh
Path=
Terminal=False
Type=Application
